# ansible/roles/ddcci/files/ddcci-nvidia.rules
ACTION=="add", SUBSYSTEM=="i2c", ATTR{name}=="NVIDIA*", RUN+="/bin/sh -c 'echo ddcci 0x37 > /sys$devpath/new_device'"
